# crontab
* * * * * /root/oto/crashdetect.sh
1,31 * * * * /root/oto/kill.sh
