# crontab
* * * * * /root/oto/crashdetect.sh
1 1,9,17 * * * /root/oto/kill.sh
